# root/etc/init.d/warp
#!/bin/sh /etc/rc.common
# Cloudflare WARP service for OpenWrt using usque.
# SPDX-License-Identifier: GPL-3.0-or-later

START=99
STOP=10
USE_PROCD=1

USQUE_BIN="/usr/bin/usque"
CONFIG_FILE="/etc/warp/config.json"
SOCKS_PID_FILE="/var/run/warp-socks.pid"
RUNTIME_DIR="/var/run/warp"
RUNTIME_CONFIG_FILE="$RUNTIME_DIR/config.json"
LOG_FILE="/var/log/warp.log"
NFT_TABLE="luci_warp"

warp_log() {
	mkdir -p /var/log 2>/dev/null || true
	printf '%s warp: %s\n' "$(date '+%Y-%m-%d %H:%M:%S')" "$*" >> "$LOG_FILE" 2>/dev/null || true
	logger -t warp "$*" 2>/dev/null || true
}

find_usque_iface() {
	if [ -s "$RUNTIME_DIR/tun_iface" ]; then
		cat "$RUNTIME_DIR/tun_iface"
		return
	fi

	ip link show 2>/dev/null | awk -F': ' '/^[0-9]+: tun[0-9]+/ {print $2; exit}'
}

usque_pids() {
	if command -v pgrep >/dev/null 2>&1; then
		pgrep -f "usque.*nativetun" 2>/dev/null && return
	fi

	ps w 2>/dev/null | awk '/[u]sque/ && /nativetun/ {print $1}'
}

ensure_tun_available() {
	modprobe tun >/dev/null 2>&1 || true

	if [ ! -c /dev/net/tun ]; then
		warp_log "Missing /dev/net/tun. Install kmod-tun or use firmware with TUN built in."
		return 1
	fi
}

force_stop_usque() {
	pids="$(usque_pids)"
	[ -n "$pids" ] || return 0

	kill $pids 2>/dev/null || true
	count=0
	while [ "$count" -lt 5 ]; do
		sleep 1
		pids="$(usque_pids)"
		[ -z "$pids" ] && return 0
		count=$((count + 1))
	done

	warp_log "Force killing stuck usque process: $pids"
	kill -9 $pids 2>/dev/null || true
}

find_wan_default_route() {
	tun_iface="$1"
	ip -4 route show default | awk -v tun="$tun_iface" '
		{
			gw = "";
			dev = "";
			for (i = 1; i <= NF; i++) {
				if ($i == "via") gw = $(i + 1);
				if ($i == "dev") dev = $(i + 1);
			}
			if (dev != "" && dev != tun) {
				print gw "|" dev;
				exit;
			}
		}'
}

resolve_endpoint_ip() {
	endpoint_host="$1"
	resolve_host_ip "$endpoint_host"
}

resolve_host_ip() {
	endpoint_host="$1"
	endpoint_ip="$(nslookup "$endpoint_host" 2>/dev/null | awk '
		/^Address[[:space:]]+[0-9]+:/ && $3 ~ /^[0-9.]+$/ {print $3}
		/^Address[0-9]*:/ && $2 ~ /^[0-9.]+$/ {print $2}
	' | tail -n 1)"
	[ -z "$endpoint_ip" ] && endpoint_ip="$(resolveip -4 "$endpoint_host" 2>/dev/null | head -n 1)"
	printf '%s\n' "$endpoint_ip"
}

is_endpoint_route_ip() {
	printf '%s\n' "$1" | awk -F. '
		NF != 4 { exit 1 }
		{
			for (i = 1; i <= 4; i++) {
				if ($i !~ /^[0-9]+$/ || $i < 0 || $i > 255) exit 1
			}

			if ($1 == 198 && ($2 == 18 || $2 == 19)) exit 1
			exit 0
		}'
}

uci_endpoint_ips() {
	config_get endpoint_h2_v4 config endpoint_h2_v4 ""
	case "$endpoint_h2_v4" in
		*.*.*.*)
			is_endpoint_route_ip "$endpoint_h2_v4" && printf '%s\n' "$endpoint_h2_v4"
			;;
	esac
}

config_endpoint_ips() {
	if command -v jsonfilter >/dev/null 2>&1 && [ -f "$CONFIG_FILE" ]; then
		for expr in '@.endpoint_v4' '@.endpoint_h2_v4'; do
			endpoint_ip="$(jsonfilter -i "$CONFIG_FILE" -e "$expr" 2>/dev/null)"
			case "$endpoint_ip" in
				*.*.*.*)
					is_endpoint_route_ip "$endpoint_ip" && printf '%s\n' "$endpoint_ip"
					;;
			esac
		done
	fi
}

endpoint_ips() {
	endpoint_host="$1"
	{
		uci_endpoint_ips
		config_endpoint_ips
		resolve_host_ip "$endpoint_host"
	} | awk 'NF && !seen[$0]++'
}

json_safe_value() {
	case "$1" in
		""|*[!A-Za-z0-9.:_-]*) return 1 ;;
	esac
}

json_set_string() {
	json_file="$1"
	json_key="$2"
	json_value="$3"
	tmp_file="$RUNTIME_DIR/config.$$"

	[ -s "$json_file" ] || return 1
	json_safe_value "$json_value" || return 1

	if grep -q "\"$json_key\"[[:space:]]*:" "$json_file"; then
		sed "s/\"$json_key\"[[:space:]]*:[[:space:]]*\"[^\"]*\"/\"$json_key\": \"$json_value\"/" "$json_file" > "$tmp_file" && \
			mv "$tmp_file" "$json_file"
		return
	fi

	awk -v key="$json_key" -v value="$json_value" '
		!inserted {
			pos = index($0, "{")
			if (pos > 0) {
				print substr($0, 1, pos)
				print "  \"" key "\": \"" value "\","
				rest = substr($0, pos + 1)
				if (rest != "") print rest
				inserted = 1
				next
			}
		}
		{ print }
		END { if (!inserted) exit 1 }
	' "$json_file" > "$tmp_file" && mv "$tmp_file" "$json_file"
}

prepare_runtime_config() {
	config_get ENDPOINT_H2_V4 config endpoint_h2_v4 ""

	rm -f "$RUNTIME_CONFIG_FILE"
	if [ -z "$ENDPOINT_H2_V4" ]; then
		printf '%s\n' "$CONFIG_FILE"
		return 0
	fi

	case "$ENDPOINT_H2_V4" in
		*.*.*.*) ;;
		*)
			warp_log "Invalid endpoint_h2_v4: $ENDPOINT_H2_V4"
			return 1
			;;
	esac

	if ! is_endpoint_route_ip "$ENDPOINT_H2_V4"; then
		warp_log "Invalid endpoint_h2_v4: $ENDPOINT_H2_V4 is reserved or fake-ip"
		return 1
	fi

	cp "$CONFIG_FILE" "$RUNTIME_CONFIG_FILE" 2>/dev/null || {
		warp_log "Failed to create runtime usque config"
		return 1
	}

	json_set_string "$RUNTIME_CONFIG_FILE" endpoint_h2_v4 "$ENDPOINT_H2_V4" || {
		warp_log "Failed to set endpoint_h2_v4 in runtime usque config"
		rm -f "$RUNTIME_CONFIG_FILE"
		return 1
	}

	warp_log "Using HTTP/2 endpoint override: $ENDPOINT_H2_V4"
	printf '%s\n' "$RUNTIME_CONFIG_FILE"
}

pre_proxy_url() {
	config_get PRE_PROXY_ENABLED config pre_proxy_enabled "0"
	config_get PRE_PROXY_TYPE config pre_proxy_type "socks5"
	config_get PRE_PROXY_ADDR config pre_proxy_addr "127.0.0.1"
	config_get PRE_PROXY_PORT config pre_proxy_port "7890"

	[ "$PRE_PROXY_ENABLED" = "1" ] || return 1
	[ -n "$PRE_PROXY_ADDR" ] && [ -n "$PRE_PROXY_PORT" ] || return 1

	case "$PRE_PROXY_TYPE" in
		http|socks5) ;;
		*)
			warp_log "Unsupported pre-proxy type: $PRE_PROXY_TYPE"
			return 1
			;;
	esac

	case "$PRE_PROXY_ADDR" in
		*:*) proxy_addr="[$PRE_PROXY_ADDR]" ;;
		*) proxy_addr="$PRE_PROXY_ADDR" ;;
	esac

	printf '%s://%s:%s\n' "$PRE_PROXY_TYPE" "$proxy_addr" "$PRE_PROXY_PORT"
}

protect_endpoint_routes() {
	endpoint_host="$1"
	wan_gw="$2"
	wan_if="$3"

	[ -n "$wan_if" ] || return 0

	for end_ip in $(endpoint_ips "$endpoint_host"); do
		case "$end_ip" in
			*.*.*.*) ;;
			*) continue ;;
		esac
		if ! is_endpoint_route_ip "$end_ip"; then
			warp_log "Ignored reserved or fake endpoint route: $end_ip"
			continue
		fi

		if [ -n "$wan_gw" ]; then
			if ip route replace "$end_ip/32" via "$wan_gw" dev "$wan_if" 2>/dev/null; then
				warp_log "Added direct route to endpoint $end_ip via $wan_gw dev $wan_if"
			else
				warp_log "Failed to add direct route to endpoint $end_ip via $wan_gw dev $wan_if"
			fi
		else
			if ip route replace "$end_ip/32" dev "$wan_if" 2>/dev/null; then
				warp_log "Added direct route to endpoint $end_ip dev $wan_if"
			else
				warp_log "Failed to add direct route to endpoint $end_ip dev $wan_if"
			fi
		fi
	done
}

delete_endpoint_routes() {
	endpoint_host="$1"

	for end_ip in $(endpoint_ips "$endpoint_host"); do
		case "$end_ip" in
			*.*.*.*) ;;
			*) continue ;;
		esac
		is_endpoint_route_ip "$end_ip" || continue
		ip route del "$end_ip/32" >/dev/null 2>&1 || true
	done
}

setup_masquerade() {
	tun_iface="$1"

	if command -v nft >/dev/null 2>&1; then
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
		if nft add table inet "$NFT_TABLE" >/dev/null 2>&1 && \
			nft add chain inet "$NFT_TABLE" postrouting '{ type nat hook postrouting priority 100; policy accept; }' >/dev/null 2>&1 && \
			nft add rule inet "$NFT_TABLE" postrouting meta nfproto ipv4 oifname "$tun_iface" masquerade >/dev/null 2>&1; then
			warp_log "Installed nftables masquerade for $tun_iface"
			return 0
		fi
		warp_log "Failed to install nftables masquerade rule"
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
	fi

	if command -v iptables >/dev/null 2>&1; then
		iptables -t nat -C POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || \
			iptables -t nat -A POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || true
		warp_log "Installed iptables masquerade for $tun_iface"
		return 0
	fi

	warp_log "No nftables or iptables command found; LAN forwarding may not be masqueraded"
}

teardown_masquerade() {
	tun_iface="$1"

	if command -v nft >/dev/null 2>&1; then
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
	fi

	if [ -n "$tun_iface" ] && command -v iptables >/dev/null 2>&1; then
		iptables -t nat -D POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || true
	fi
}

configure_runtime() {
	config_load warp
	config_get ADDRESS_V4 config address_v4 ""
	config_get ADDRESS_V6 config address_v6 ""
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"
	config_get GLOBAL_PROXY config global_proxy "0"
	config_get SOCKS_ENABLED config socks_enabled "1"
	config_get SOCKS_PORT config socks_port "1080"
	config_get BYPASS_CHINA config bypass_china "0"

	max_wait=30
	count=0
	USQ_IF=""

	while [ "$count" -lt "$max_wait" ]; do
		USQ_IF="$(find_usque_iface)"
		[ -n "$USQ_IF" ] && break
		sleep 1
		count=$((count + 1))
	done

	if [ -z "$USQ_IF" ]; then
		warp_log "Failed to find usque TUN interface after ${max_wait}s"
		exit 1
	fi

	mkdir -p "$RUNTIME_DIR"
	echo "$USQ_IF" > "$RUNTIME_DIR/tun_iface"
	warp_log "usque interface created: $USQ_IF"

	[ -n "$ADDRESS_V4" ] && ip addr add "$ADDRESS_V4/32" dev "$USQ_IF" 2>/dev/null || true
	[ -n "$ADDRESS_V6" ] && ip -6 addr add "$ADDRESS_V6/128" dev "$USQ_IF" 2>/dev/null || true
	ip link set "$USQ_IF" up 2>/dev/null || true

	end_host="${ENDPOINT%:*}"
	wan_route="$(find_wan_default_route "$USQ_IF")"
	wan_gw="${wan_route%%|*}"
	wan_if="${wan_route#*|}"
	protect_endpoint_routes "$end_host" "$wan_gw" "$wan_if"

	if [ "$GLOBAL_PROXY" = "1" ]; then
		ip route replace 0.0.0.0/1 dev "$USQ_IF" 2>/dev/null || true
		ip route replace 128.0.0.0/1 dev "$USQ_IF" 2>/dev/null || true
	fi

	setup_masquerade "$USQ_IF"

	if [ "$BYPASS_CHINA" = "1" ] && [ -x /usr/bin/warp-update-china ]; then
		/usr/bin/warp-update-china apply >/dev/null 2>&1 || true
	fi

	if [ "$SOCKS_ENABLED" = "1" ]; then
		if command -v microsocks >/dev/null 2>&1; then
			if [ -n "$ADDRESS_V4" ]; then
				microsocks -i 0.0.0.0 -p "$SOCKS_PORT" -b "$ADDRESS_V4" &
			else
				microsocks -i 0.0.0.0 -p "$SOCKS_PORT" &
			fi
			echo $! > "$SOCKS_PID_FILE"
			warp_log "Started microsocks on port $SOCKS_PORT"
		else
			warp_log "microsocks is not installed; SOCKS5 proxy was not started"
		fi
	fi
}

optimize_sysctl() {
	# 增加 UDP 接收和发送缓冲区大小，优化 QUIC 性能
	sysctl -w net.core.rmem_max=2621440 >/dev/null 2>&1 || true
	sysctl -w net.core.wmem_max=2621440 >/dev/null 2>&1 || true
	warp_log "Optimized system network parameters for QUIC"
}

start_service() {
	optimize_sysctl
	config_load warp
	config_get ENABLED config enabled 0
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"
	config_get HTTP2 config http2 "0"
	config_get SNI_ADDRESS config sni_address ""
	[ "$ENABLED" = "1" ] || return 0

	if [ ! -x "$USQUE_BIN" ]; then
		warp_log "Missing usque binary at $USQUE_BIN"
		return 1
	fi

	if [ ! -f "$CONFIG_FILE" ]; then
		warp_log "Missing $CONFIG_FILE. Run warp-manager register first."
		return 1
	fi

	ensure_tun_available || return 1
	mkdir -p "$RUNTIME_DIR"
	ACTIVE_CONFIG="$(prepare_runtime_config)" || return 1
	PRE_PROXY_URL="$(pre_proxy_url)"

	end_host="${ENDPOINT%:*}"
	wan_route="$(find_wan_default_route "")"
	wan_gw="${wan_route%%|*}"
	wan_if="${wan_route#*|}"
	protect_endpoint_routes "$end_host" "$wan_gw" "$wan_if"

	procd_open_instance "usque"
	procd_set_param command "$USQUE_BIN" -c "$ACTIVE_CONFIG" nativetun
	[ "$HTTP2" = "1" ] && procd_append_param command --http2
	[ -n "$SNI_ADDRESS" ] && procd_append_param command -s "$SNI_ADDRESS"
	if [ -n "$PRE_PROXY_URL" ]; then
		procd_set_param env HTTP_PROXY="$PRE_PROXY_URL" HTTPS_PROXY="$PRE_PROXY_URL" ALL_PROXY="$PRE_PROXY_URL" NO_PROXY="127.0.0.1,localhost,::1"
	fi
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param term_timeout 5
	procd_close_instance

	if [ "$HTTP2" = "1" ]; then
		warp_log "Started usque nativetun using HTTP/2"
	else
		warp_log "Started usque nativetun using HTTP/3"
	fi
	if [ -n "$PRE_PROXY_URL" ]; then
		if [ "$HTTP2" = "1" ]; then
			warp_log "Using pre-proxy for HTTP/2: $PRE_PROXY_URL"
		else
			warp_log "Pre-proxy is configured but usque only honors proxy environment in HTTP/2 mode"
		fi
	fi
	configure_runtime &
}

stop_service() {
	config_load warp
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"

	USQ_IF="$(find_usque_iface)"

	if [ -f "$SOCKS_PID_FILE" ]; then
		kill "$(cat "$SOCKS_PID_FILE")" 2>/dev/null || true
		rm -f "$SOCKS_PID_FILE"
	fi

	if [ -n "$USQ_IF" ]; then
		ip route del 0.0.0.0/1 dev "$USQ_IF" >/dev/null 2>&1 || true
		ip route del 128.0.0.0/1 dev "$USQ_IF" >/dev/null 2>&1 || true
		teardown_masquerade "$USQ_IF"
	fi

	end_host="${ENDPOINT%:*}"
	delete_endpoint_routes "$end_host"

	[ -x /usr/bin/warp-update-china ] && /usr/bin/warp-update-china clear >/dev/null 2>&1 || true
	force_stop_usque
	rm -f "$RUNTIME_DIR/tun_iface" "$RUNTIME_CONFIG_FILE"
	warp_log "Service stopped"
}

reload_service() {
	stop_service
	start_service
}
