# root/etc/init.d/warp
#!/bin/sh /etc/rc.common
# Cloudflare WARP service for OpenWrt using usque.
# SPDX-License-Identifier: GPL-3.0-or-later

START=99
STOP=10
USE_PROCD=1

USQUE_BIN="/usr/bin/usque"
CONFIG_FILE="/etc/warp/config.json"
SOCKS_PID_FILE="/var/run/warp-socks.pid"
RUNTIME_DIR="/var/run/warp"
LOG_FILE="/var/log/warp.log"
NFT_TABLE="luci_warp"

warp_log() {
	mkdir -p /var/log 2>/dev/null || true
	printf '%s warp: %s\n' "$(date '+%Y-%m-%d %H:%M:%S')" "$*" >> "$LOG_FILE" 2>/dev/null || true
	logger -t warp "$*" 2>/dev/null || true
}

find_usque_iface() {
	if [ -s "$RUNTIME_DIR/tun_iface" ]; then
		cat "$RUNTIME_DIR/tun_iface"
		return
	fi

	ip link show 2>/dev/null | awk -F': ' '/^[0-9]+: tun[0-9]+/ {print $2; exit}'
}

usque_pids() {
	if command -v pgrep >/dev/null 2>&1; then
		pgrep -f "usque.*nativetun" 2>/dev/null && return
	fi

	ps w 2>/dev/null | awk '/[u]sque/ && /nativetun/ {print $1}'
}

ensure_tun_available() {
	modprobe tun >/dev/null 2>&1 || true

	if [ ! -c /dev/net/tun ]; then
		warp_log "Missing /dev/net/tun. Install kmod-tun or use firmware with TUN built in."
		return 1
	fi
}

force_stop_usque() {
	pids="$(usque_pids)"
	[ -n "$pids" ] || return 0

	kill $pids 2>/dev/null || true
	count=0
	while [ "$count" -lt 5 ]; do
		sleep 1
		pids="$(usque_pids)"
		[ -z "$pids" ] && return 0
		count=$((count + 1))
	done

	warp_log "Force killing stuck usque process: $pids"
	kill -9 $pids 2>/dev/null || true
}

find_wan_default_route() {
	tun_iface="$1"
	ip route show default | awk -v tun="$tun_iface" '
		{
			gw = "";
			dev = "";
			for (i = 1; i <= NF; i++) {
				if ($i == "via") gw = $(i + 1);
				if ($i == "dev") dev = $(i + 1);
			}
			if (dev != "" && dev != tun) {
				print gw "|" dev;
				exit;
			}
		}'
}

resolve_endpoint_ip() {
	endpoint_host="$1"
	resolve_host_ip "$endpoint_host"
}

resolve_host_ip() {
	endpoint_host="$1"
	endpoint_ip="$(nslookup "$endpoint_host" 2>/dev/null | awk '/^Address[[:space:]]+[0-9]+:/ {print $3} /^Address[0-9]*:/ {print $2}' | tail -n 1)"
	[ -z "$endpoint_ip" ] && endpoint_ip="$(resolveip -4 "$endpoint_host" 2>/dev/null | head -n 1)"
	printf '%s\n' "$endpoint_ip"
}

config_endpoint_ips() {
	if command -v jsonfilter >/dev/null 2>&1 && [ -f "$CONFIG_FILE" ]; then
		for expr in '@.endpoint_v4' '@.endpoint_h2_v4'; do
			endpoint_ip="$(jsonfilter -i "$CONFIG_FILE" -e "$expr" 2>/dev/null)"
			case "$endpoint_ip" in
				*.*.*.*) printf '%s\n' "$endpoint_ip" ;;
			esac
		done
	fi
}

endpoint_ips() {
	endpoint_host="$1"
	{
		config_endpoint_ips
		resolve_host_ip "$endpoint_host"
	} | awk 'NF && !seen[$0]++'
}

protect_endpoint_routes() {
	endpoint_host="$1"
	wan_gw="$2"
	wan_if="$3"

	[ -n "$wan_if" ] || return 0

	for end_ip in $(endpoint_ips "$endpoint_host"); do
		if [ -n "$wan_gw" ]; then
			ip route replace "$end_ip/32" via "$wan_gw" dev "$wan_if" 2>/dev/null || true
			warp_log "Added direct route to endpoint $end_ip via $wan_gw dev $wan_if"
		else
			ip route replace "$end_ip/32" dev "$wan_if" 2>/dev/null || true
			warp_log "Added direct route to endpoint $end_ip dev $wan_if"
		fi
	done
}

delete_endpoint_routes() {
	endpoint_host="$1"

	for end_ip in $(endpoint_ips "$endpoint_host"); do
		ip route del "$end_ip/32" >/dev/null 2>&1 || true
	done
}

setup_masquerade() {
	tun_iface="$1"

	if command -v nft >/dev/null 2>&1; then
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
		if nft add table inet "$NFT_TABLE" >/dev/null 2>&1 && \
			nft add chain inet "$NFT_TABLE" postrouting '{ type nat hook postrouting priority 100; policy accept; }' >/dev/null 2>&1 && \
			nft add rule inet "$NFT_TABLE" postrouting meta nfproto ipv4 oifname "$tun_iface" masquerade >/dev/null 2>&1; then
			warp_log "Installed nftables masquerade for $tun_iface"
			return 0
		fi
		warp_log "Failed to install nftables masquerade rule"
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
	fi

	if command -v iptables >/dev/null 2>&1; then
		iptables -t nat -C POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || \
			iptables -t nat -A POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || true
		warp_log "Installed iptables masquerade for $tun_iface"
		return 0
	fi

	warp_log "No nftables or iptables command found; LAN forwarding may not be masqueraded"
}

teardown_masquerade() {
	tun_iface="$1"

	if command -v nft >/dev/null 2>&1; then
		nft delete table inet "$NFT_TABLE" >/dev/null 2>&1 || true
	fi

	if [ -n "$tun_iface" ] && command -v iptables >/dev/null 2>&1; then
		iptables -t nat -D POSTROUTING -o "$tun_iface" -j MASQUERADE >/dev/null 2>&1 || true
	fi
}

configure_runtime() {
	config_load warp
	config_get ADDRESS_V4 config address_v4 ""
	config_get ADDRESS_V6 config address_v6 ""
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"
	config_get GLOBAL_PROXY config global_proxy "0"
	config_get SOCKS_ENABLED config socks_enabled "1"
	config_get SOCKS_PORT config socks_port "1080"
	config_get BYPASS_CHINA config bypass_china "0"

	max_wait=30
	count=0
	USQ_IF=""

	while [ "$count" -lt "$max_wait" ]; do
		USQ_IF="$(find_usque_iface)"
		[ -n "$USQ_IF" ] && break
		sleep 1
		count=$((count + 1))
	done

	if [ -z "$USQ_IF" ]; then
		warp_log "Failed to find usque TUN interface after ${max_wait}s"
		exit 1
	fi

	mkdir -p "$RUNTIME_DIR"
	echo "$USQ_IF" > "$RUNTIME_DIR/tun_iface"
	warp_log "usque interface created: $USQ_IF"

	[ -n "$ADDRESS_V4" ] && ip addr add "$ADDRESS_V4/32" dev "$USQ_IF" 2>/dev/null || true
	[ -n "$ADDRESS_V6" ] && ip -6 addr add "$ADDRESS_V6/128" dev "$USQ_IF" 2>/dev/null || true
	ip link set "$USQ_IF" up 2>/dev/null || true

	end_host="${ENDPOINT%:*}"
	wan_route="$(find_wan_default_route "$USQ_IF")"
	wan_gw="${wan_route%%|*}"
	wan_if="${wan_route#*|}"
	protect_endpoint_routes "$end_host" "$wan_gw" "$wan_if"

	if [ "$GLOBAL_PROXY" = "1" ]; then
		ip route replace 0.0.0.0/1 dev "$USQ_IF" 2>/dev/null || true
		ip route replace 128.0.0.0/1 dev "$USQ_IF" 2>/dev/null || true
	fi

	setup_masquerade "$USQ_IF"

	if [ "$BYPASS_CHINA" = "1" ] && [ -x /usr/bin/warp-update-china ]; then
		/usr/bin/warp-update-china apply >/dev/null 2>&1 || true
	fi

	if [ "$SOCKS_ENABLED" = "1" ]; then
		if command -v microsocks >/dev/null 2>&1; then
			if [ -n "$ADDRESS_V4" ]; then
				microsocks -i 0.0.0.0 -p "$SOCKS_PORT" -b "$ADDRESS_V4" &
			else
				microsocks -i 0.0.0.0 -p "$SOCKS_PORT" &
			fi
			echo $! > "$SOCKS_PID_FILE"
			warp_log "Started microsocks on port $SOCKS_PORT"
		else
			warp_log "microsocks is not installed; SOCKS5 proxy was not started"
		fi
	fi
}

start_service() {
	config_load warp
	config_get ENABLED config enabled 0
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"
	config_get HTTP2 config http2 "0"
	config_get SNI_ADDRESS config sni_address ""
	[ "$ENABLED" = "1" ] || return 0

	if [ ! -x "$USQUE_BIN" ]; then
		warp_log "Missing usque binary at $USQUE_BIN"
		return 1
	fi

	if [ ! -f "$CONFIG_FILE" ]; then
		warp_log "Missing $CONFIG_FILE. Run warp-manager register first."
		return 1
	fi

	ensure_tun_available || return 1
	mkdir -p "$RUNTIME_DIR"

	end_host="${ENDPOINT%:*}"
	wan_route="$(find_wan_default_route "")"
	wan_gw="${wan_route%%|*}"
	wan_if="${wan_route#*|}"
	protect_endpoint_routes "$end_host" "$wan_gw" "$wan_if"

	procd_open_instance "usque"
	procd_set_param command "$USQUE_BIN" -c "$CONFIG_FILE" nativetun
	[ "$HTTP2" = "1" ] && procd_append_param command --http2
	[ -n "$SNI_ADDRESS" ] && procd_append_param command -s "$SNI_ADDRESS"
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param term_timeout 5
	procd_close_instance

	if [ "$HTTP2" = "1" ]; then
		warp_log "Started usque nativetun using HTTP/2"
	else
		warp_log "Started usque nativetun using HTTP/3"
	fi
	configure_runtime &
}

stop_service() {
	config_load warp
	config_get ENDPOINT config endpoint "engage.cloudflareclient.com:2408"

	USQ_IF="$(find_usque_iface)"

	if [ -f "$SOCKS_PID_FILE" ]; then
		kill "$(cat "$SOCKS_PID_FILE")" 2>/dev/null || true
		rm -f "$SOCKS_PID_FILE"
	fi

	if [ -n "$USQ_IF" ]; then
		ip route del 0.0.0.0/1 dev "$USQ_IF" >/dev/null 2>&1 || true
		ip route del 128.0.0.0/1 dev "$USQ_IF" >/dev/null 2>&1 || true
		teardown_masquerade "$USQ_IF"
	fi

	end_host="${ENDPOINT%:*}"
	delete_endpoint_routes "$end_host"

	[ -x /usr/bin/warp-update-china ] && /usr/bin/warp-update-china clear >/dev/null 2>&1 || true
	force_stop_usque
	rm -f "$RUNTIME_DIR/tun_iface"
	warp_log "Service stopped"
}

reload_service() {
	stop_service
	start_service
}
